# database/baseX/supervisord.conf
[program:java]
command = /home/dotcloud/current/run.sh
directory = /home/dotcloud/current
redirect_stderr = true
